# .nemesis/install/com.nemesis.publisher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template do LaunchAgent do nemesis-publisher.
     O placeholder {{EXEC}} e substituido por `nemesis-publisher --install-service`.
     launchd nao tem equivalente de EnvironmentFile: DATABASE_URL e afins sao lidos
     de {{ENV_FILE}} pelo proprio publisher no startup (fallback in-process). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nemesis.publisher</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{EXEC}}</string>
        <string>--serve</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
